# config/job-listing.service
[Unit]
Description=Job Listing Application - Go Web Server
Documentation=https://github.com/LionelPinheiroDuarte/job-listing
After=network.target
# After=network.target means: start only after network is ready

[Service]
Type=simple
# Type=simple: the process runs in foreground

User=appuser
Group=appuser
# User/Group: run as non-root user for security

WorkingDirectory=/opt/job-listing
# WorkingDirectory: where the binary and data.json are located

Environment="PORT=8000"
Environment="GIN_MODE=release"
# Environment: set environment variables for the app

ExecStart=/opt/job-listing/main
# ExecStart: command to start the application (the compiled Go binary)

Restart=always
# Restart=always: if the app crashes, systemd automatically restarts it

RestartSec=10
# RestartSec=10: wait 10 seconds before restarting

StandardOutput=journal
StandardError=journal
# StandardOutput/Error: send logs to journald (view with journalctl)

[Install]
WantedBy=multi-user.target
# WantedBy=multi-user.target: start this service at normal system boot


# COMPARISON WITH DOCKER
# 
# Docker equivalent:
#   restart: always           → Restart=always
#   ports: "8000:8000"        → Not needed in systemd (direct port access)
#   environment: PORT=8000    → Environment="PORT=8000"
#   user: appuser             → User=appuser
#   healthcheck:              → No built-in equivalent (need custom script)
#
# What systemd does that Docker doesn't:
#   - Integrated with system boot (WantedBy)
#   - Automatic logging to journald
#   - Process supervision at OS level
#
# What Docker does that systemd doesn't:
#   - Dependency isolation (libraries, etc.)
#   - Port mapping abstraction
#   - Built-in health checks